# examples/fss-server.service.in
[Unit]
Description=Flight-Safety-System Server
Requires=postgresql.service
After=postgresql.service

[Service]
ExecStart=@PREFIX@/fss-server

[Install]
WantedBy=multi-user.target
