# examples/fss-server.service.in
[Unit]
Description=Flight-Safety-System Server

[Service]
ExecStart=@PREFIX@/fss-server

[Install]
WantedBy=multi-user.target
